# pcmanfm/stuff/autostart/pcmanfm.desktop
[Desktop Entry]
Type=Application
Name=PCmanFM Desktop Manager
Name[fr]=Gestionnaire de bureau PCmanFM
Exec=pcmanfm -d
NotShowIn=XFCE;Razor;
