# volumeicon/stuff/autostart/volumeicon.desktop
[Desktop Entry]
Type=Application
Name=Audio Volume mamangement
Name[fr]=Gestion du volume audio
Exec=volumeicon
NotShowIn=XFCE;Razor;
