# lxpanel/stuff/autostart/lxpanel.desktop
[Desktop Entry]
Type=Application
Name=LXPanel Desktop panel
Name[fr]=Panneau du bureau LXPanel
Exec=lxpanel -p slitaz
NotShowIn=XFCE;Razor;
